# platform/centec/sonic-platform-modules-ieu/debian/platform-modules-ieu-ctc8180.postinst
systemctl enable ieu_ctc8180_platform.service
systemctl start ieu_ctc8180_platform.service
