# platform/centec-arm64/sonic-platform-modules-e530/debian/platform-modules-e530-24x2q.postinst
systemctl enable 24x2q_platform.service
systemctl start 24x2q_platform.service
sed -i "/    docker create --privileged/a\        -v /usr/lib/python2.7/dist-packages/sonic_platform:/usr/lib/python2.7/dist-packages/sonic_platform:rw \\\\" /usr/bin/pmon.sh
